# build/supervisor-app.conf
[program:papertrail-shipper]
command = rsyslogd -n -f /src/build/papertrail.conf

[program:scheduler]
command = bash /src/docker-run.sh